# private/benchmarks/benchmark-agent-automation/templates/etc/systemd/system/teamcity-agent.service
[Unit]
Description=TeamCity Agent

[Service]
LimitNOFILE=262140
LimitCORE=infinity
User=teamcity
ExecStart=/opt/teamcity-agent/bin/agent.sh run

[Install]
WantedBy=multi-user.target
